# ClaudeCode-Trade/deploy/aws/systemd/claudecode-api.service
[Unit]
Description=ClaudeCode-Trade API (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claudebot
Group=claudebot
WorkingDirectory=/opt/claudecode-trade
EnvironmentFile=/opt/claudecode-trade/.env
ExecStart=/usr/bin/node src/api/server.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claudecode-api

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/opt/claudecode-trade/data /var/log/claudecode-trade
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
