# ClaudeCode-Trade/deploy/aws/systemd/claudecode-backup.timer
[Unit]
Description=Run ClaudeCode-Trade DB backup daily at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=claudecode-backup.service

[Install]
WantedBy=timers.target
